# deploy/systemd/chrononode-anchor@.timer
[Unit]
Description=Weekly ChronoNode Arweave Checkpoint Anchor Timer for %i
Documentation=https://github.com/chrononode/chrononode

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
